# ntfy/app/supervisord.conf
[supervisord]
nodaemon=true
user=root

[program:app]
command=/etc/ntfy/entrypoint.sh
autostart=true
autorestart=true
